# templates/etc/profile.d/python.sh.j2
export PATH={{ python_configure_prefix }}/bin:$PATH
